# service/supervisord.conf
[program:service]
command=/setup/service-start.sh
numprocs=1
autostart=true
autorestart=true
